# Mailbot/crontab
0 2 * * * cd /app && dotnet Mailbot.dll >> /var/log/mailbot.log 2>&1
