# EmailSync/crontab
0 2 * * * cd /app && dotnet ApplicationTracker.EmailSync.dll >> /var/log/emailsync.log 2>&1
